# supervisor.webvirtmgr.conf
[program:webvirtmgr]
command=/bin/sh -c "/webvirtmgr/bootstrap.sh && /usr/bin/python /webvirtmgr/manage.py run_gunicorn -c /webvirtmgr/conf/gunicorn.conf.py"
directory=/webvirtmgr
autostart=true
autorestart=true
stdout_logfile=/var/log/webvirtmgr.log
stderr_logfile=/var/log/webvirtmgr.log
user=webvirtmgr

[program:webvirtmgr-console]
command=/usr/bin/python /webvirtmgr/console/webvirtmgr-console -v -d
directory=/webvirtmgr
autostart=true
autorestart=true
stdout_logfile=/var/log/webvirtmgr-console.log
stderr_logfile=/var/log/webvirtmgr-console.log
user=webvirtmgr
